# carm.service
[Unit]
Description=CaRM

[Service]
Type=simple
ExecStart=/etc/carm/carm /etc/carm/carm.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
